# supervisor.conf
[program:boobiesbot]
command=/usr/local/bin/launch-boobiesbot-container
directory=/
user=root
